# .config/systemd/user/kmonad@.service
[Unit]
Description=KMonad daemon

[Service]
Type=exec
ExecStart=kmonad "%h/.config/kmonad/%i.conf"

[Install]
WantedBy=default.target
